# systemd/nexus_ai.service
[Unit]
Description=NEXUS AI Trading Signal System
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/nexus_ai
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/ubuntu/nexus_ai/.env
ExecStart=/home/ubuntu/nexus_ai/venv/bin/python main.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nexus_ai

# Resource limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
